# data/session-migration.conf
description "Session Migration"
author "Iain Lane <iain.lane@canonical.com>"

start on startup

task

exec session-migration
